# supervisor/conf.d/socat.conf
[program:socat]
command=/usr/bin/socat tcp6-listen:8880,fork tcp6:[::1]:8080
stopsignal=KILL
autostart=true
autorestart=true
stdout_logfile=/var/log/socat.log
stderr_logfile=/var/log/socat.log
